# api/deploy/vibecheckapi.service
[Unit]
Description=Vibecheck API Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=deploy
Group=www-data
WorkingDirectory=/var/www/vibecheck/api
ExecStart=/var/www/vibecheck/api/target/release/vibecheckapi --serve --port 4000 --database /var/www/vibecheck/api/vibecheckapi.db
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
