# api/deploy/vibecheckapi.service
[Unit]
Description=Vibecheck API Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=deploy
Group=www-data
WorkingDirectory=/var/www/vibecheckapi
ExecStart=/var/www/vibecheckapi/target/release/vibecheckapi --serve --port 3000 --database /var/www/vibecheckapi/vibecheckapi.db
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
