# api/deploy/vibecheckapi.service
[Unit]
Description=yo API Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=deploy
Group=www-data
WorkingDirectory=/var/www/yo/api
ExecStart=/var/www/yo/api/target/release/yoapi --serve --port 4000 --database /var/www/yo/api/yoapi.db
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
